# qqqq.conf
[program:qqqq]
command=python3 /home/curator2/qqqq/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/qqqq/test.err.log
stdout_logfile=/home/curator2/qqqq/test.out.log